# fluctlight-serve.service
[Unit]
Description=FluctlightDB in-process HTTP API
After=network.target

[Service]
Type=simple
User=fluctlight
Group=fluctlight
EnvironmentFile=/etc/fluctlight/environment
Environment=FLUCTLIGHT_STORAGE=v4
Environment=FLUCTLIGHT_RATE_CAPACITY=120
Environment=FLUCTLIGHT_RATE_REFILL=60
ExecStart=/bin/sh -ec 'exec "$FLUCTLIGHT_INSTALL/target/release/fluctlight" serve --addr 127.0.0.1:8792 --path "$FLUCTLIGHT_BRAIN_PATH"'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
