# systemd/fluctlight-drill.timer
[Unit]
Description=Fluctlight brain integrity drill

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
