# systemd/fluctlight-embed.service
[Unit]
Description=Fluctlight embedding sidecar (optional)
After=network.target

[Service]
Type=simple
User=fluctlight
Group=fluctlight
EnvironmentFile=/etc/fluctlight/environment
ExecStart=/bin/sh -ec 'cd "$FLUCTLIGHT_INSTALL/embed-server" && exec .venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8793'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
