# systemd/fluctlight-drill.service
[Unit]
Description=Fluctlight weekly verify drill

[Service]
Type=oneshot
User=fluctlight
Group=fluctlight
EnvironmentFile=/etc/fluctlight/environment
ExecStart=/bin/sh -ec 'exec "$FLUCTLIGHT_INSTALL/scripts/fluctlight-drill.sh"'
