# systemd/fluctlight-backup.service
[Unit]
Description=Fluctlight brain backup

[Service]
Type=oneshot
EnvironmentFile=/etc/fluctlight/environment
ExecStart=/bin/sh -ec 'exec "$FLUCTLIGHT_INSTALL/scripts/fluctlight-backup.sh"'
