# systemd/serverbrain-bridge.service.d/fluctlight.conf
[Service]
Environment=FLUCTLIGHT_SERVE_URL=http://127.0.0.1:8792
Environment=FLUCTLIGHT_EMBED_URL=http://127.0.0.1:8793/embed
Environment=FLUCTLIGHT_BRAIN=/home/ambugo/.fluctlight/tenants/default/brain
Environment=FLUCTLIGHT_SDK_PATH=/home/ambugo/fluctlightdb/sdks/python
EnvironmentFile=-/home/ambugo/.fluctlight/auth.env
# auth.env may define FLUCTLIGHT_API_KEY for bridge (same secret as serve keys entry)
