# systemd/fluctlight-stream.service.d/fluctlight.conf
[Service]
Environment=FLUCTLIGHT_SERVE_URL=http://127.0.0.1:8792
Environment=FLUCTLIGHT_EMBED_URL=http://127.0.0.1:8793/embed
