# systemd/fluctlight-serve.service.d/production.conf
[Service]
Environment=FLUCTLIGHT_STORAGE=v4
Environment=FLUCTLIGHT_RATE_CAPACITY=120
Environment=FLUCTLIGHT_RATE_REFILL=60
Environment=FLUCTLIGHT_CHECKPOINT_EVERY_N=50
Environment=FLUCTLIGHT_MAX_HOT_TENANTS=256
Environment=FLUCTLIGHT_MAX_CONNECTIONS=500
Environment=FLUCTLIGHT_WAL_FSYNC=batched
Environment=FLUCTLIGHT_WAL_BATCH_MS=25
Environment=FLUCTLIGHT_WAL_BATCH_RECORDS=32
EnvironmentFile=-/etc/fluctlight/auth.env
# Agent clients use FLUCTLIGHT_API_KEY matching a key from auth.env
